# Library/LaunchAgents/clash.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>clash</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
		<string>Background</string>
		<string>LoginWindow</string>
		<string>StandardIO</string>
		<string>System</string>
	</array>
	<key>SoftResourceLimits</key>
	<dict>
		<key>OpenFiles</key>
		<integer>1048576</integer>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/etc/profiles/per-user/{{ .chezmoi.username }}/bin/clash</string>
		<string>-f</string>
		<string>{{ .chezmoi.homeDir }}/.config/clash/config.yaml</string>
		<string>-d</string>
		<string>{{ .chezmoi.homeDir }}/.local/state/clash</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/clash.log</string>
	<key>StandardOutPath</key>
	<string>{{ .chezmoi.homeDir }}/Library/Logs/clash.log</string>
</dict>
</plist>
